# infra/roles/web/templates/imdbgraph-api.service
[Unit]
Description=IMDB Graph api service

[Service]
User=clubhouse-bot
ExecStart=/home/imdbgraph-api/start.sh

[Install]
WantedBy=multi-user.target